# deploy/pi-gen/stage-clawrig/04-configure-watchdog/files/clawrig-gateway-watchdog.timer
[Unit]
Description=Check OpenClaw gateway health every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
RandomizedDelaySec=15

[Install]
WantedBy=timers.target
